# celery-queue/configs/supervisord.conf

[supervisord]
nodaemon=true

[supervisorctl]


[program:workers]
numprocs = 10
process_name = queue_%(process_num)s
command=celery -A tasks worker -l info -E -Q queue_%(process_num)s -c 1